# supervisor-config/supervisor-nginx.conf
[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
