# supervisor-config/supervisor-php5.conf
[program:php5]
command=/usr/sbin/php5-fpm
numprocs=1
autostart=true
autorestart=true
